# settings/utapau/crontab
# Vergeet niet te updaten na editen:
#   crontab ~/crontab

# m h  dom mon dow   command
@reboot echo "\* Doe import-irclogs" >> ~/vergeet-mij-niet

